# lib/systemd/system/fwkit.service
[Unit]
Description=SD57 Firewall Kit Control Service
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/default/fwkit
ExecStart=/usr/sbin/fwkit rules load
ExecReload=/usr/sbin/fwkit rules reload
ExecStop=/usr/sbin/fwkit rules unload

[Install]
WantedBy=multi-user.target
